# snapserver.socket
[Socket]
ListenFIFO=%t/snapserver/fifo

[Install]
WantedBy=sockets.target
